# ansible/roles/dynhostname/files/dynhostname.cron
# Detect changes to our public IP (e.g. association with an elastic IP) and
# update the custom DNS name.
SHELL=/bin/bash

* * * * * root /usr/local/bin/check-hostname-consistency || systemctl reload dynhostname
