# ansible/roles/tcpd_update/templates/tcpd-update.cron.j2
8,38 * * * * root /usr/local/sbin/tcpd-update -d http://datastore.{{ subd }}
