# ansible/roles/tcpd_update/templates/tcpd-update.cron.j2
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
8,38 * * * * root /usr/local/sbin/tcpd-update -d http://datastore.{{ subd }}
