# ansible/roles/firewall_blocklist/templates/fwbl-update.cron.j2
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
*/5 * * * * root /usr/local/sbin/fwbl-update -d datastore.{{ subd }} -r {{ aws_layout[aws_status][aws_profile_s3].region }}
