# ansible/roles/kitbag/templates/kitbag.cron.j2
#!/bin/bash

{{ comint_base }}/dispatch/kitbag
