# ansible/roles/carbon/files/carbon.cron
#!/bin/bash

find /var/lib/carbon/whisper/ -type f -mtime +14 -delete
find /var/lib/carbon/whisper/ -empty -delete
