# ansible/roles/kitbag/files/kitbag.cron
#!/bin/bash

/usr/local/collaboration/dispatch/kitbag
