# ansible/roles/rbldnsd_sync/files/dnsbl-sync.cron
 * * * * * root /usr/local/sbin/dnsbl-sync
