# ansible/roles/rbldnsd_sync/files/rbl-sync.cron
 * * * * * root /usr/local/sbin/rbl-sync
