# ansible/roles/base/files/dailyjobs.cron
# Run the daily, weekly, and monthly jobs
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/local/sbin:/usr/local/bin
MAILTO=root

# run-parts
00 0 * * * root run-parts /etc/cron.daily
30 0 * * 0 root run-parts /etc/cron.weekly
00 1 1 * * root run-parts /etc/cron.monthly
